# despliegue/supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn Restaurantes.wsgi:application --bind 0.0.0.0:8000
directory=/restaurantes/manage.py
user=root
autostart=true
autorestart=true
redirect_stderr=true
